# gx19.service
[Unit]
Description=gx19 daemon service

[Service]
ExecStart=/usr/bin/gxdaemon

[Install]
WantedBy=default.target